# chef/cookbooks/class2go-base-ubuntu/files/default/update-mnt-perms.conf
description     "Make sure mnt has correct perms"

start on startup
task

script
    if [ ! -d /mnt ]; then
        mkdir /mnt
    fi
    chmod 1777 /mnt
end script

