# chef/cookbooks/class2go-ops-dns/files/default/update-route53-dns.conf
description     "Create DNS record for this host in route53"

start on startup
task

exec /usr/sbin/update-route53-dns.sh
