# chef/cookbooks/class2go-ops-dns/files/default/update-route53-dns.conf
# update-route53-dns - set up DNS for this machine

description     "Create DNS entries for this hostname"

start on static-netowrk-up

task
exec /usr/sbin/update-route53-dns.sh
